# config/supervisor/sshd.conf
[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true
stdout_events_enabled=true
stderr_events_enabled=true
