# config/supervisor/x2go.conf
[program:x2go]
command=/usr/sbin/x2gocleansessions
autostart=true
autorestart=false
stdout_events_enabled=true
stderr_events_enabled=true
